# home/Library/LaunchAgents/localhost.blackd.plist.tmpl
{{- if eq .chezmoi.os "darwin" }}
{{- if stat (joinPath .chezmoi.homeDir ".rye/shims/blackd") -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>localhost.blackd</string>
        <key>Program</key>
        <string>{{- (joinPath .chezmoi.homeDir ".rye/shims/blackd") -}}</string>
        <key>RunAtLoad</key>
        <true />
        <key>KeepAlive</key>
        <true />
    </dict>
</plist>
{{- end -}}
{{- end -}}
